# cmd/dashboard/docker/dashboard.sv.conf
[program:nginx]
command = /usr/sbin/nginx -g "daemon off;"
autorestart = true

[program:dashboard]
command = /usr/local/bin/dashboard --docker-key-dir /etc/nuclio/dashboard/registry-credentials --listen-addr :18070
autorestart = true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
