# vendor/gopkg.in/jcmturner/gokrb5.v7/test/testenv/older-krb5kdc-vagrant/krb5kdc.service
[Unit]
Description=Kerberos 5 KDC
After=syslog.target network.target

[Service]
Type=forking
ExecStart=/usr/local/sbin/krb5kdc 

[Install]
WantedBy=multi-user.target
