# vendor/gopkg.in/jcmturner/gokrb5.v7/test/testenv/mit-krb5kdc/krb5kdc-shorttickets.service
[Unit]
Description=Kerberos KDC
After=docker.service
After=network.target
Requires=docker.service

[Service]
Environment="DOCKER_IMAGE=jcmturner/gokrb5:kdc-shorttickets"
Environment="PORT=58"
Environment="NAME=gokrb5-kdc-shorttickets"
TimeoutStartSec=0
#Restart=always
ExecStartPre=-/usr/bin/docker kill %n
ExecStartPre=-/usr/bin/docker rm %n
ExecStartPre=-/usr/bin/docker pull ${DOCKER_IMAGE}
ExecStart=/usr/bin/docker run -h kdc.test.gokrb5 -v /etc/localtime:/etc/localtime:ro -p ${PORT}:88 -p ${PORT}:88/udp --rm --name ${NAME} ${DOCKER_IMAGE}
ExecStartPost=/bin/sh -c "while [ ! -f ${DATA_HOST}/.initialised ]; do sleep 2; done"
ExecStop=/usr/bin/docker stop --time=60 %n
ExecStopPost=-/usr/bin/docker rm %n

[Install]
WantedBy=multi-user.target