# vendor/golang.org/x/build/vcs-test/vcweb.service
[Unit]
Description=vcs-test web server
Requires=vcweb-http.socket vcweb-https.socket

[Service]
Type=notify
ExecStart=/usr/bin/vcweb -d /home/vcweb
WatchdogSec=30s
Restart=always
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=vcweb
User=uucp

[Install]
WantedBy=multi-user.target
