# vendor/golang.org/x/build/vcs-test/vcweb-http.socket
[Unit]
Description=version control server HTTP socket
PartOf=vcweb.service

[Socket]
ListenStream=80
Service=vcweb.service

[Install]
WantedBy=sockets.target
