# vendor/golang.org/x/build/vcs-test/vcweb-https.socket
[Unit]
Description=version control server HTTPS socket
PartOf=vcweb.service

[Socket]
ListenStream=443
Service=vcweb.service

[Install]
WantedBy=sockets.target
