# vendor/golang.org/x/build/env/linux-arm64/linaro/rundockerbuildlet.service
[Unit]
Description=Run Buildlets in Docker
After=network.target

[Install]
WantedBy=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/rundockerbuildlet -basename=linaro -image=gobuilder-arm64-linaro:1 -n=4
Restart=always
RestartSec=2
StartLimitInterval=0
