# vendor/golang.org/x/build/env/linux-arm-qemu/buildlet.service
[Unit]
Description=Buildlet

[Service]
ExecStart=/usr/local/bin/stage0

[Install]
WantedBy=multi-user.target
