# vendor/golang.org/x/build/env/linux-arm-qemu/qemu.service
[Unit]
Description=QEMU

[Service]
ExecStart=/usr/local/bin/buildlet-qemu

[Install]
WantedBy=multi-user.target
