# vendor/golang.org/x/build/env/linux-arm/scaleway/rundockerbuildlet.service
[Unit]
Description=Run Buildlets in Docker
After=network.target

[Install]
WantedBy=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/rundockerbuildlet -basename=scaleway -image=gobuilder-arm-scaleway:latest -n=1 
Restart=always
RestartSec=2
StartLimitInterval=0